# temp/test_run_manager/bin/parabuild
#!/bin/sh
# parabuild     Starts/Stops Parabuild manager
#
# chkconfig: 345 99 10
# description:	The Parabuild manager provides an ability to configure
#               and run automatic and scheduled builds. Automatic builds
#               start at every check in. Scheduled builds start at
#               configurable time.
# processname:	parabuild
# pidfile:	/var/run/parabuild.pid

# Configuration variables prequisites:
#
#   PARABUILD_HOME         Must point at your Parabuild installation directory.
#                          Default installation directory is /opt/parabuild.
PARABUILD_HOME=/opt/parabuild

# Prepare Parabuild command
PARABUILD_COMMAND=$PARABUILD_HOME/bin/parabuild.sh
[ -f $PARABUILD_COMMAND ] || exit 2

RETVAL=0
PARABUILD_PID=/var/run/parabuild.pid

umask 077

# Start Parabuild
start() {
        echo -n $"Starting Parabuild: "
        su - parabuild -c "$PARABUILD_COMMAND start"
        echo
        return $RETVAL
}

# Stop Parabuild
stop() {
        echo -n $"Shutting down Parabuild: "
        su - parabuild -c "$PARABUILD_COMMAND stop"
        echo
        return $RETVAL
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  *)
        echo $"Usage: $0 {start|stop}"
        exit 0
esac

exit $?
